# deploy/systemd/media-edge.service
[Unit]
Description=VOS-RS Media Edge
After=network.target

[Service]
Type=simple
User=vos-rs
Group=vos-rs
Environment=VOS_RS_CONFIG_FILE=/etc/vos-rs/config.yaml
ExecStart=/opt/vos-rs/bin/media-edge
Restart=on-failure
RestartSec=5
LimitNOFILE=1048576
WorkingDirectory=/opt/vos-rs

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/vos-rs/data/recordings /var/log/vos-rs
PrivateTmp=true

[Install]
WantedBy=multi-user.target
